# mqtt_logger/install/systemd/maps-logger-compress.timer
[Unit]
Description=Hourly compression of MAPS Logger files

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
