# install/systemd/maps-logger.service
[Unit]
Description=MAPS Value Logger
Documentation=https://docs.mapsmessaging.io
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=/etc/maps-logger/env
ExecStart=java -jar /usr/lib/maps-logger/maps-value-logger.jar
Restart=on-failure
RestartSec=10s
StandardOutput=journal
StandardError=journal
SyslogIdentifier=maps-logger
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true

[Install]
WantedBy=multi-user.target
